# build/opentsdb/hbase_supervisor.conf
[program:hbase]
command=/opt/hbase/bin/start-hbase.sh
autorestart=true
autostart=true
startsecs=5
environment=COMPRESSION="NONE",JAVA_HOME="/usr/lib/jvm/jre",HBASE_HOME=/opt/hbase

; logging
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_logfile=log/%(program_name)s.log
